# meta-roslynos-distro/recipes-core/systemd-conf/files/70-ap-interface.rules
SUBSYSTEM=="net", KERNEL=="wlan*", ACTION=="add", RUN+="/usr/sbin/iw dev %k interface add ap%n type __ap"